# test/traffic-sim/traffic-sim.service
[Unit]
Description=LogHarbor traffic simulator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/traffic-sim
EnvironmentFile=/root/traffic-sim/.env
ExecStart=/usr/bin/python3 /root/traffic-sim/traffic-sim.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
